# daemons/app.service
After=network.target

[Service]
ExecStart=/usr/local/bin/app.sh
Type=oneshot
RemainAfterExit=true

[Install]
WantedBy=default.target
