# daemons/app.service
After=network.target

[Service]
ExecStart=/usr/local/bin/app.sh
RemainAfterExit=true
Restart=always
[Install]
WantedBy=default.target
