# daemons/server.service
After=network.target

[Service]
ExecStart=/usr/local/bin/server.sh
RemainAfterExit=true
Restart=always
[Install]
WantedBy=default.target
